# ansible/roles/systemd/files/udev/50-monitor-hotplug.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/usr/bin/systemctl start toggle-monitor.service"
