# ansible/roles/systemd/files/docker-system-prune.service
[Unit]
Description=docker system prune

[Service]
ExecStart=/usr/bin/docker system prune -f --volumes
