# ansible/roles/systemd/files/docker-system-prune.timer
[Unit]
Description=Run docker system prune daily

[Timer]
OnBootSec=1hour
OnUnitActiveSec=1day

[Install]
WantedBy=timers.target
