# ansible/roles/systemd/files/fetch-dotfiles.timer
[Unit]
Description=fetch dotfiles git repo

[Timer]
OnBootSec=1minute
OnUnitActiveSec=1hour

[Install]
WantedBy=timers.target
