# ansible/roles/systemd/files/mpc-reload.service
[Unit]
Description=Reload mpc database

[Service]
User=user
ExecStart=/bin/sh -c "/home/user/.local/bin/mpc-reload"
