# ansible/roles/systemd/files/mpc-reload.timer
[Unit]
Description=Reload mpc database

[Timer]
OnBootSec=1minute

[Install]
WantedBy=timers.target
