# ansible/roles/systemd/files/newsboat-reload.timer
[Unit]
Description=Reload newsboat feeds

[Timer]
OnBootSec=10minutes
OnUnitActiveSec=12hours

[Install]
WantedBy=timers.target
