# ansible/roles/systemd/files/pacman-update-files-database.timer
[Unit]
Description=Sync the files database

[Timer]
OnBootSec=1hour
OnUnitActiveSec=1day

[Install]
WantedBy=timers.target
